# packaging/systemd/vps-sentinel.service
[Unit]
Description=VPS Sentinel Security Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=@BIN_PATH@ daemon --config @CONFIG_PATH@
Restart=always
RestartSec=5
User=root
Group=root
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=false
ReadWritePaths=@DATA_DIR@ @LOG_DIR@
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
